# ekshell@.service
[Unit]
Description=Encrypted Kerberos-aware Rshell Server
After=local-fs.target

[Service]
ExecStart=/usr/kerberos/sbin/kshd -e
StandardInput=socket
