# gssftp.socket
[Unit]
Description=Kerberos-aware FTP Server Activation Socket

[Socket]
ListenStream=21
Accept=true

[Install]
WantedBy=sockets.target
