# klogin@.service
[Unit]
Description=Kerberos-aware Rlogin Server
After=local-fs.target

[Service]
ExecStart=/usr/kerberos/sbin/klogind
StandardInput=socket
