# krb5-telnet@.service
[Unit]
Description=Kerberos-aware Telnet Server
After=local-fs.target

[Service]
ExecStart=/usr/kerberos/sbin/telnetd
StandardInput=socket
