# krb5-telnet.socket
[Unit]
Description=Kerberos-aware Telnet Server Activation Socket

[Socket]
ListenStream=23
Accept=true

[Install]
WantedBy=sockets.target
