# kshell.socket
[Unit]
Description=Kerberos-aware Rshell Server Activation Socket

[Socket]
ListenStream=544
Accept=true

[Install]
WantedBy=sockets.target
